# service/locust.service
[Unit]
Description=Locust

[Service]
ExecStart=/root/locust/startup-locust.sh start

[Install]
WantedBy=multi-user.target